# packages/openhydroponics/openhydroponics-0.6.0-py3-none-any.whl/openhydroponics/conf/systemd/openhydroponicsd.service
[Unit]
Description=OpenHydroponics daemon

[Service]
ExecStart={path}
Restart=on-failure
User={user}

[Install]
WantedBy=default.target
